# root/etc/periodic/15min/owncloud
#!/bin/sh
sudo -u owncloud php -f /opt/owncloud/cron.php > /dev/stdout 2>&1
